# python/JobManager/src/pyjob_run.conf
description "local client to run jobs for the JobManager server"
author "FAC"

setuid fac
setgid fac

expect fork

start on runlevel [2345]
stop  on runlevel [016]

task

exec /home/fac_files/code/python/JobManager/src/pyjob_run.py

console log
